# dot-matrix-weather-fax/systemd/weather-fax.service
[Unit]
Description=Generate and print the daily weather fax
Wants=network-online.target cups.service
After=network-online.target cups.service

[Service]
Type=oneshot
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME
Environment=HOME=/home/YOUR_USERNAME
Environment=LANG=C.UTF-8
ExecStart=/usr/local/bin/weather-fax --print
